# backups/phase2_20260625_170311/com.bigshot.receive-payment.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.bigshot.receive-payment</string>

  <key>ProgramArguments</key>
  <array>
    <string>/usr/bin/python3</string>
    <string>/Users/bigshot/ai-automation/business-ai/scripts/receive_payment_server.py</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/bigshot/ai-automation/business-ai</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>RECEIVE_PAYMENT_HOST</key>
    <string>127.0.0.1</string>
    <key>RECEIVE_PAYMENT_PORT</key>
    <string>5059</string>
  </dict>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/Users/bigshot/ai-automation/business-ai/logs/receive_payment_server.out.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/bigshot/ai-automation/business-ai/logs/receive_payment_server.err.log</string>
</dict>
</plist>
